# templates/debian/prerm
#!/bin/sh
set -e
APP_NAME="${APP_NAME}"
case "$1" in
    remove|purge)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop "${APP_NAME}" 2>/dev/null || true
            systemctl disable "${APP_NAME}" 2>/dev/null || true
            rm -f /etc/systemd/system/${SERVICE_FILE} || :
            systemctl daemon-reload || :
        fi
        killall -9 "${APP_NAME}" 2>/dev/null || true
        rm -f "/usr/local/bin/${APP_USER}" || :
        rm -f "/usr/share/applications/${DESKTOP_FILE}" || :
        rm -f "/usr/share/pixmaps/${ICON_FILE}" || :
        ;;
esac
exit 0
